# crontabs
00 06 * * * php [ROOT]/admin/watchdog/watchdog.php